# selfhosted-moodle/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=/usr/local/sbin/php-fpm
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true

[program:moodle-cron]
command=/bin/bash -c "while true; do php /var/www/moodle/admin/cli/cron.php --keep-alive=0 >/dev/null 2>&1; sleep 60; done"
autostart=true
autorestart=true
